# site/minecraft/templates/minecraft.service.epp
[Unit]
Description=Mineraft server

Wants=network.target
After=network.targer

[Service]
WorkingDirectory=<%= $install_dir %>
ExecStart=/usr/bin/java -Xms512M -Xms32M -jar server.jar nogui

[Install]
WantedBy=multi-user.target
